# deploy/systemd/jih-autodeploy.timer
[Unit]
Description=Poll GitHub for new commits every minute and auto-deploy Java Interview Hub

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
